# deploy/systemd/trading_bot.service
[Unit]
Description=Trading Bot Worker
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/Garima
ExecStart=/home/ubuntu/Garima/.venv/bin/python src/bot/main.py
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/Garima/.env

[Install]
WantedBy=multi-user.target